# fs/etc/services.d/jc-agent/run
#!/usr/bin/with-contenv bash

echo "Starting agent..."

/opt/jc/bin/jumpcloud-agent -conf=/opt/jc/jcagent.conf